# packaging/allinone/s6/rspamd/run
#!/bin/sh
# rspamd refuses to run as root; drop privileges to the _rspamd service user
exec /command/s6-setuidgid _rspamd /usr/bin/rspamd --no-fork -c /etc/rspamd/rspamd.conf
