# packaging/allinone/s6/clamd/run
#!/bin/sh
exec /command/s6-setuidgid clamav /usr/sbin/clamd --foreground
